# memcached/1.4.10/supervisord.conf

[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true

[program:memcache]
command=/bin/bash -c "/opt/soft/memcache-1.4.10-runtime/bin/memcached -d  -p 10036 -u root"
autorestart=true

